# packages/ginggar/ginggar-3.0.1554-py3-none-any.whl/ginggar/_meta/misc/crontab
*/5 * * * *  root  /usr/bin/curl --insecure https://localhost/ginggar/api/feeds/actions/crawl/ &> /dev/null
